# gen_iso/profile/airootfs/etc/systemd/system/pacman-init.service
[Unit]
Description=Initialize Pacman keyring and optionally trust ArchZFS key
Requires=etc-pacman.d-gnupg.mount
After=etc-pacman.d-gnupg.mount time-sync.target
BindsTo=etc-pacman.d-gnupg.mount
Before=archlinux-keyring-wkd-sync.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Initialize and populate keyring only; no DKMS actions here
ExecStart=/usr/bin/pacman-key --init
ExecStart=/usr/bin/pacman-key --populate
# Optionally fetch and locally sign the ArchZFS key to allow archzfs repo access if needed
ExecStart=/bin/sh -lc 'set -e; key=F75D9D76AE2D87A471F0DCD02B4CF6B7663BDAAE; pacman-key -r "$key" || true; pacman-key --lsign-key "$key" || true'

[Install]
WantedBy=multi-user.target
# Also ensure zfs-dkms-autoinstall.service gets pulled in
Alias=zfs-dkms-autoinstall.service
